# deployment/arena-watcher.service
[Unit]
Description=Arena Watcher Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=michaelsik12
WorkingDirectory=/home/michaelsik12/arena-watcher/current
EnvironmentFile=/etc/arena-watcher/arena-watcher.env
ExecStart=/usr/bin/dotnet /home/michaelsik12/arena-watcher/current/DiscordBot.dll
Restart=always
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=30
SyslogIdentifier=arena-watcher

[Install]
WantedBy=multi-user.target
